# systemd/animaksm-swap-proxy.service
[Unit]
Description=Animaksm Deduplicating Swap Proxy (experimental)
Documentation=https://github.com/animaios/animaksm
BindsTo=animaksm.service
After=animaksm.service

[Service]
Type=simple
ExecStart=/usr/local/bin/animaksm-swap-proxy run --store-path /var/lib/animaksm/pagestore.dat
ExecStopPost=-/usr/local/bin/animaksm-swap-proxy cleanup
Restart=no

# Security hardening
AmbientCapabilities=CAP_SYS_ADMIN CAP_SYS_RAWIO CAP_IPC_LOCK
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYS_RAWIO CAP_IPC_LOCK
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/animaksm

[Install]
WantedBy=multi-user.target
